# wayland/.config/systemd/user/wlsunset.service.d/options.conf
[Service]
ExecStart=wlsunset -l 48.8 -L 2.2
