# miscellaneous/.config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=batsignal -d 5 -c 10 -w 20
